# contrib/systemd/xdg-desktop-portal-termfilechooser.service.in
[Unit]
Description=Portal service (terminal file chooser implementation)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.termfilechooser
ExecStart=@libexecdir@/xdg-desktop-portal-termfilechooser
Restart=on-failure
